# data/watchful.service
[Unit]
Description=Run watchful service

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /opt/watchful/monitor.py
